# claude/launchd/com.skanehira.utility-self-improving.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skanehira.utility-self-improving</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>set -o pipefail; cd /Users/skanehira/dev/github.com/skanehira/dotfiles &amp;&amp; perl -e 'alarm shift @ARGV; exec @ARGV' 3600 claude -p "/utility-self-improving" --permission-mode acceptEdits --allowedTools "Bash,Read,Edit,Write,Glob,Grep,Agent" --output-format json | tee /Users/skanehira/.claude/logs/self-improving.json | jq -r '.[] | select(.type=="result") | .result'; rc=$?; if [ $rc -ne 0 ]; then terminal-notifier -title "utility-self-improving" -message "週次自己改善が失敗 (exit $rc)" -subtitle "~/.claude/logs/self-improving.err を確認"; fi; exit $rc</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/skanehira/.claude/logs/self-improving.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/skanehira/.claude/logs/self-improving.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/skanehira/.local/bin:/Users/skanehira/.deno/bin:/Users/skanehira/.nix-profile/bin:/etc/profiles/per-user/skanehira/bin:/run/current-system/sw/bin:/nix/var/nix/profiles/default/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/skanehira</string>
        <key>GHQ_ROOT</key>
        <string>/Users/skanehira/dev</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
